# deploy/cms/neurau-cms.service
[Unit]
Description=Neurau CMS OAuth proxy and webhook
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/neurau.eu
EnvironmentFile=/etc/neurau-cms.env
ExecStart=/usr/local/bin/neurau-cms
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/neurau.eu
PrivateTmp=true

[Install]
WantedBy=multi-user.target
